# scripts/linux/eduky-monitor.service
[Unit]
Description=EDUKY-商品监控系统
After=network.target
Wants=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/eduky-monitor
Environment="PYTHONPATH=/opt/eduky-monitor"
Environment="FLASK_ENV=production"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/opt/eduky-monitor/venv/bin/python /opt/eduky-monitor/app_production.py
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=eduky-monitor

# 安全设置
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/opt/eduky-monitor /var/log/eduky-monitor
ProtectHome=yes

[Install]
WantedBy=multi-user.target